# service/pipewire/run
#!/bin/sh
exec 2>&1
XDG_RUNTIME_DIR=/run/user/1000 exec pipewire 2>/dev/null
